# etc/systemd/aegyrad.service
[Unit]
Description=Aegyra face-auth daemon
Documentation=https://github.com/archledger/aegyra
After=systemd-tmpfiles-setup.service
RequiresMountsFor=/run /etc/aegyra

[Service]
Type=simple
ExecStart=/usr/local/bin/aegyrad
Restart=on-failure
RestartSec=2s

# Needs root: TPM (/dev/tpmrm0), camera (/dev/video*), /run/aegyra.sock,
# and read access to /etc/aegyra for sealed envelope + models.
User=root
Group=root

# Hardening. Keep CAP_SYS_ADMIN off; we only need DAC/device access.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run /etc/aegyra
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Devices: TPM resource manager + V4L cameras.
DeviceAllow=/dev/tpmrm0 rw
DeviceAllow=char-video4linux rw
PrivateDevices=no

# The socket is created by the daemon itself at /run/aegyra.sock (chmod 0666
# with per-op uid enforcement). A future revision will switch to 0660 with an
# `aegyra` group and socket activation.

[Install]
WantedBy=multi-user.target
