# etc/systemd/linhellod.service
[Unit]
Description=LinuxHello face-auth daemon
Documentation=https://github.com/archledger/linhello
After=systemd-tmpfiles-setup.service
RequiresMountsFor=/run /etc/linhello

[Service]
Type=simple
ExecStart=/usr/local/bin/linhellod
Restart=on-failure
RestartSec=2s

# Fail-closed liveness: if /etc/linhello/antispoof.onnx is missing, refuse to
# verify rather than silently falling back to the device-trust-only gate.
# Override in a drop-in (/etc/systemd/system/linhellod.service.d/) if you need
# to run without the anti-spoof model during early bring-up.
Environment=LINHELLO_REQUIRE_ANTISPOOF=1

# Needs root: TPM (/dev/tpmrm0), camera (/dev/video*), /run/linhello.sock,
# and read access to /etc/linhello for sealed envelope + models.
User=root
Group=root

# Hardening. Drop every capability except the few the daemon actually uses:
#   CAP_CHOWN          — chown the socket to the `linhello` group at startup
#   CAP_DAC_OVERRIDE   — read/write sealed envelopes + open the TPM/V4L device
#                        nodes regardless of their owner/group mode bits
#   CAP_FOWNER         — manage modes on files it creates under /etc/linhello
# Without an explicit bounding set the daemon would keep the full root cap set
# (CAP_SYS_ADMIN, CAP_SYS_MODULE, CAP_NET_ADMIN, CAP_SYS_PTRACE, …) despite this
# section's intent. AmbientCapabilities stays empty (it is already root).
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run /etc/linhello
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# @resources is needed by onnxruntime's thread pool (sched_setaffinity);
# keep @privileged stripped so bpf/kexec/init_module stay denied.
SystemCallFilter=~@privileged
# ...but the daemon must chown its socket to the `linhello` group at startup,
# and @privileged strips the chown family. Re-allow just that. Without this the
# daemon dies SIGSYS on boot the moment the linhello group exists.
SystemCallFilter=@chown

# Devices: TPM resource manager + V4L cameras.
DeviceAllow=/dev/tpmrm0 rw
DeviceAllow=char-video4linux rw
PrivateDevices=no

# The socket is created by the daemon itself at /run/linhello.sock (chmod 0666
# with per-op uid enforcement). A future revision will switch to 0660 with an
# `linhello` group and socket activation.

[Install]
WantedBy=multi-user.target
